# files/launch/audiofocusmanager.conf.in
start on lsm-ready

respawn

script
        exec @WEBOS_INSTALL_SBINDIR@/audiofocusmanager
end script
